# cue/systemd/cue.timer
[Unit]
Description=Cue daily release notification timer

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target
