# hwaiting/systemd/hwaiting-backup.service
[Unit]
Description=Back up hwaiting SQLite database via restic
After=network-online.target
Wants=network-online.target
OnFailure=herald@%n

[Service]
Type=oneshot
User=ajussi

LoadCredentialEncrypted=hwaiting-backup-url:/etc/credstore/hwaiting/hwaiting-backup-url.cred
LoadCredentialEncrypted=hwaiting-backup-pass:/etc/credstore/hwaiting/hwaiting-backup-pass.cred

Environment=RESTIC_REPOSITORY_FILE=%d/hwaiting-backup-url
Environment=RESTIC_PASSWORD_FILE=%d/hwaiting-backup-pass

ExecStart=/usr/local/sbin/hwaiting-backup.sh

PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6
